# deployment/atlas_os.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.markus.atlasos.runtime</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>runtime.atlas_runtime_daemon</string>
    <string>--interval</string>
    <string>60</string>
    <string>--llm-model</string>
    <string>gpt5.5</string>
    <string>--db-path</string>
    <string>runtime/state/atlas_runtime.sqlite</string>
    <string>--log-path</string>
    <string>runtime/logs/atlas_runtime.log</string>
    <string>--inbox-dir</string>
    <string>runtime/events/inbox</string>
    <string>--ui-inbox-path</string>
    <string>runtime/inbox/user_event.jsonl</string>
    <string>--market-config-path</string>
    <string>runtime/config/user_config.json</string>
    <string>--proactive-update-every-seconds</string>
    <string>7200</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/markus/AtlasOS</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/markus/AtlasOS/runtime/logs/atlas_daemon.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/markus/AtlasOS/runtime/logs/atlas_daemon.err.log</string>
</dict>
</plist>
